# packaging/linux/systemd/loramapr-receiverd.service
[Unit]
Description=LoRaMapr Receiver Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=loramapr
Group=loramapr
ExecStart=/usr/bin/loramapr-receiverd -config /etc/loramapr/receiver.json
WorkingDirectory=/var/lib/loramapr
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
