# packaging/debian/scripts/postinst
#!/bin/sh
set -eu

create_service_user() {
  if ! getent group loramapr >/dev/null 2>&1; then
    addgroup --system loramapr >/dev/null 2>&1 || true
  fi

  if ! getent passwd loramapr >/dev/null 2>&1; then
    adduser \
      --system \
      --home /var/lib/loramapr \
      --no-create-home \
      --ingroup loramapr \
      --disabled-login \
      --disabled-password \
      --shell /usr/sbin/nologin \
      loramapr >/dev/null 2>&1 || true
  fi
}

ensure_serial_access() {
  if ! getent passwd loramapr >/dev/null 2>&1; then
    return
  fi

  if ! getent group dialout >/dev/null 2>&1; then
    echo "dialout group not found; serial access may require manual host setup." >&2
    return
  fi

  if id -nG loramapr 2>/dev/null | tr ' ' '\n' | grep -qx dialout; then
    return
  fi

  if command -v usermod >/dev/null 2>&1; then
    usermod -a -G dialout loramapr >/dev/null 2>&1 || adduser loramapr dialout >/dev/null 2>&1 || true
  else
    adduser loramapr dialout >/dev/null 2>&1 || true
  fi
}

ensure_runtime_permissions() {
  mkdir -p /etc/loramapr /var/lib/loramapr /var/log/loramapr
  chmod 0755 /etc/loramapr || true
  chown -R loramapr:loramapr /var/lib/loramapr /var/log/loramapr || true
  chmod 0750 /var/lib/loramapr /var/log/loramapr || true

  if [ -f /var/lib/loramapr/receiver-state.json ]; then
    chown loramapr:loramapr /var/lib/loramapr/receiver-state.json || true
    chmod 0600 /var/lib/loramapr/receiver-state.json || true
  fi

  if [ -f /etc/loramapr/receiver.json ]; then
    chown root:loramapr /etc/loramapr/receiver.json || true
    chmod 0640 /etc/loramapr/receiver.json || true
  fi
}

normalize_legacy_systemd_unit() {
  legacy_unit="/etc/systemd/system/loramapr-receiverd.service"
  packaged_unit="/lib/systemd/system/loramapr-receiverd.service"
  if [ ! -f "${legacy_unit}" ]; then
    return
  fi
  if [ ! -f "${packaged_unit}" ]; then
    return
  fi

  if cmp -s "${legacy_unit}" "${packaged_unit}"; then
    rm -f "${legacy_unit}" || true
    return
  fi

  echo "Keeping legacy unit override at ${legacy_unit}; review if migration is required." >&2
}

start_or_restart_service() {
  if ! command -v systemctl >/dev/null 2>&1; then
    return
  fi

  systemctl daemon-reload || true
  action="${1:-install}"
  if [ "${action}" = "install" ]; then
    systemctl enable loramapr-receiverd.service >/dev/null 2>&1 || true
  fi

  if ! systemctl is-enabled --quiet loramapr-receiverd.service; then
    return
  fi

  if systemctl is-active --quiet loramapr-receiverd.service; then
    systemctl restart loramapr-receiverd.service || true
  else
    systemctl start loramapr-receiverd.service || true
  fi
}

case "${1:-}" in
  configure|reconfigure)
    create_service_user
    ensure_serial_access
    normalize_legacy_systemd_unit
    ensure_runtime_permissions

    if [ -n "${2:-}" ]; then
      start_or_restart_service "upgrade"
    else
      start_or_restart_service "install"
    fi
    ;;
esac

exit 0
